# deploy/oracle/systemd/probate-bot-sync.service
[Unit]
Description=Probate Bot daily sync
After=network.target

[Service]
Type=oneshot
EnvironmentFile=/etc/default/probate-bot
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/deploy/oracle/run-sync.sh
StandardOutput=append:__APP_DIR__/logs/daily-sync.log
StandardError=append:__APP_DIR__/logs/daily-sync.log
