# deploy/oracle/systemd/probate-bot-backup.timer
[Unit]
Description=Run Probate Bot backup every day

[Timer]
OnCalendar=*-*-* 05:45:00
Persistent=true

[Install]
WantedBy=timers.target
